# deploy/pco-unifi-sync.service
[Unit]
Description=Planning Center Calendar -> UniFi Access Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/fccplanner2unifi
EnvironmentFile=/root/fccplanner2unifi/.env
ExecStart=/root/fccplanner2unifi/.venv/bin/python -m uvicorn py_app.main:app --host 0.0.0.0 --port 3000
Restart=always
RestartSec=5

# Hardening (keep modest; controller may need network access)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
